# geoip-autoupdate.timer
[Unit]
Description=Daily GeoIP / GeoSite update

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
